# backend/deployment/mysite_supervisor.conf
[program:mysite_shahbaz] 
<<<<<<< HEAD:backend/deployment/mysite_gunicorn_supervisor.conf
command = /home/shahbazz/app/react-django/backend/deployment/gunicorn_start          ; Command to start app 
=======
command = /home/shahbazz/app/react-django/backend/gunicorn_start          ; Command to start app 
>>>>>>> 6d03c791b9fb4d0ba94f530b8774b0a88e707af3:backend/deployment/mysite_supervisor.conf
user = shahbazz                                                                                     ; User to run as 
stdout_logfile = /home/shahbazz/app/react-django/backend/logs/gunicorn_supervisor.log  ; Where to write log messages 
redirect_stderr = true                                                                        ; Save stderr in the same log 
environment=LANG=en_US.UTF-8,LC_ALL=en_US.UTF-8 
